# packaging/distribution/etc/systemd/system/tremor.service
[Unit]
Description=Tremor
Documentation=https://docs.tremor.rs
After=network-online.target
Requires=network-online.target

[Service]
User=tremor
Group=tremor
ExecStart=/bin/sh -c 'usr/bin/tremor server run --logger-config /etc/tremor/logger.yaml -f /etc/tremor/config/*.yaml -f /etc/tremor/config/*.trickle'
Restart=always
SyslogIdentifier=tremor

[Install]
WantedBy=multi-user.target
